# infra/launchd/com.lip.ollama.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent template. ``__HOME__`` is substituted with $HOME when
  ``task ollama:install`` materializes the plist into
  ``~/Library/LaunchAgents/com.lip.ollama.plist``. launchd does NOT
  expand env vars in plists, so the substitution is a build step.

  LaunchAgent ONLY. Install under ~/Library/LaunchAgents/, never under
  /Library/LaunchDaemons/. The latter loads as root and lacks the user
  context Ollama needs (model paths under ~/.ollama, GPU access, etc.).
  No <UserName> key is set deliberately so a misplacement into a system
  domain fails loudly rather than silently inheriting root identity.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lip.ollama</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ollama</string>
        <string>serve</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OLLAMA_KEEP_ALIVE</key>
        <string>300s</string>
        <key>OLLAMA_NUM_PARALLEL</key>
        <string>1</string>
        <key>OLLAMA_MAX_LOADED_MODELS</key>
        <string>1</string>
        <key>OLLAMA_FLASH_ATTENTION</key>
        <string>1</string>
        <key>OLLAMA_KV_CACHE_TYPE</key>
        <string>q8_0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/ollama/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/ollama/stderr.log</string>
</dict>
</plist>
